# bamboo-server/supervisor.conf
[program:bamboo-server]
command=/start/bamboo-server
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
autorestart=true
