# bamboo-agent/supervisor/bamboo-agent.conf
[program:bamboo-agent]
command=/start/bamboo-agent
user=bamboo
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
autorestart=true
