# bamboo-agent/supervisor.conf
[program:bamboo-agent]
command=/usr/bin/bamboo-agent
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
